# ops/crontab
# Powerfulwords cron schedule.
# Install:   crontab ops/crontab        (replaces your existing crontab)
# Inspect:   crontab -l
# Edit live: crontab -e
#
# IMPORTANT: edit the EIA_API_KEY line below to your own key
# (free at https://www.eia.gov/opendata/register.php) before installing,
# or unset it and rely on ~/.bashrc / a wrapper script.

EIA_API_KEY=__SET_ME__

# Weekly event-study rerun — Sundays 02:30
30 2 * * 0 /usr/bin/python3 /home/hongbingchen/powerfulwords/scripts/run_event_study_weekly.py >> /tmp/oil_pipeline.log 2>&1

# Oil narrative pipeline — register, fetch, prices, regimes, ingest, extract, score, hypotheses, hourly at :05
5 * * * * /usr/bin/python3 /home/hongbingchen/powerfulwords/scripts/init_sources.py >> /tmp/oil_pipeline.log 2>&1 && /usr/bin/python3 /home/hongbingchen/powerfulwords/scripts/fetch_sources.py >> /tmp/oil_pipeline.log 2>&1 && /usr/bin/python3 /home/hongbingchen/powerfulwords/scripts/fetch_prices.py --period 1mo >> /tmp/oil_pipeline.log 2>&1 && /usr/bin/python3 /home/hongbingchen/powerfulwords/scripts/compute_regimes.py >> /tmp/oil_pipeline.log 2>&1 && /usr/bin/python3 /home/hongbingchen/powerfulwords/scripts/ingest_folder.py >> /tmp/oil_pipeline.log 2>&1 && /usr/bin/python3 /home/hongbingchen/powerfulwords/scripts/extract_narratives.py >> /tmp/oil_pipeline.log 2>&1 && /usr/bin/python3 /home/hongbingchen/powerfulwords/scripts/score_narratives.py >> /tmp/oil_pipeline.log 2>&1 && /usr/bin/python3 /home/hongbingchen/powerfulwords/scripts/test_strategy_hypotheses.py >> /tmp/oil_pipeline.log 2>&1

# Composite backtest — nightly at 03:15 (after the 23:05/00:05/01:05/02:05 hourly ticks have settled latest data)
15 3 * * * /usr/bin/python3 /home/hongbingchen/powerfulwords/scripts/run_composite_backtest.py >> /tmp/oil_pipeline.log 2>&1

# Paper trading snapshot — nightly at 07:00 (after WTI 17:00 ET = 05:00 UTC+8 settlement,
# plus margin for the 06:05 hourly fetch_prices to publish the official daily close)
0 7 * * * /usr/bin/python3 /home/hongbingchen/powerfulwords/scripts/snapshot_paper_trades.py >> /tmp/oil_pipeline.log 2>&1

# AI judgment review — nightly at 07:15 (after the paper snapshot at 07:00).
# Requires ANTHROPIC_API_KEY in the crontab env to do anything (otherwise skips cleanly).
15 7 * * * /usr/bin/python3 /home/hongbingchen/powerfulwords/scripts/generate_ai_review.py >> /tmp/oil_pipeline.log 2>&1

# Daily news report (headlines + latest inventory, no scores) — nightly at 07:30
30 7 * * * /usr/bin/python3 /home/hongbingchen/powerfulwords/scripts/daily_news_report.py >> /tmp/oil_pipeline.log 2>&1
