# ops/crontab
# Powerfulwords cron schedule.
# Install:   crontab ops/crontab        (replaces your existing crontab)
# Inspect:   crontab -l
# Edit live: crontab -e
#
# IMPORTANT: edit the EIA_API_KEY line below to your own key
# (free at https://www.eia.gov/opendata/register.php) before installing,
# or unset it and rely on ~/.bashrc / a wrapper script.

EIA_API_KEY=__SET_ME__

# Weekly event-study rerun — Sundays 02:30
30 2 * * 0 /usr/bin/python3 /home/hongbingchen/powerfulwords/scripts/run_event_study_weekly.py >> /tmp/oil_pipeline.log 2>&1

# Oil narrative pipeline — hourly at :05. Wrapped in flock so that if a
# prior hour's extract_narratives is still running (large upload batches
# can take >60min), the new hour exits silently instead of clobbering
# the DB. Skipped hours catch up naturally — fetchers pull a 14d rolling
# window and the ingester is checksum-deduped.
5 * * * * flock -n /tmp/oil_pipeline.lock /home/hongbingchen/powerfulwords/scripts/hourly_pipeline.sh >> /tmp/oil_pipeline.log 2>&1

# Composite backtest — nightly at 03:15 (after the 23:05/00:05/01:05/02:05 hourly ticks have settled latest data)
15 3 * * * /usr/bin/python3 /home/hongbingchen/powerfulwords/scripts/run_composite_backtest.py >> /tmp/oil_pipeline.log 2>&1

# Technical-only backtest — nightly at 03:20 (right after composite backtest)
20 3 * * * /usr/bin/python3 /home/hongbingchen/powerfulwords/scripts/run_technical_backtest.py >> /tmp/oil_pipeline.log 2>&1

# Paper trading snapshot — nightly at 07:00 (after WTI 17:00 ET = 05:00 UTC+8 settlement,
# plus margin for the 06:05 hourly fetch_prices to publish the official daily close)
0 7 * * * /usr/bin/python3 /home/hongbingchen/powerfulwords/scripts/snapshot_paper_trades.py >> /tmp/oil_pipeline.log 2>&1

# AI judgment review — nightly at 07:15 (after the paper snapshot at 07:00).
# Requires ANTHROPIC_API_KEY in the crontab env to do anything (otherwise skips cleanly).
15 7 * * * /usr/bin/python3 /home/hongbingchen/powerfulwords/scripts/generate_ai_review.py >> /tmp/oil_pipeline.log 2>&1

# Daily news report (headlines + latest inventory, no scores) — nightly at 07:30
30 7 * * * /usr/bin/python3 /home/hongbingchen/powerfulwords/scripts/daily_news_report.py >> /tmp/oil_pipeline.log 2>&1
